# supervisord.conf
[supervisord]
nodaemon=true


[program:docsify]
command=docsify serve --port 80
autorestart=true
startretries=9999
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
